# src/modules/hotspot_nm/filesystem/root/etc/NetworkManager/dispatcher.d/10-ratos-hotspot
#!/bin/bash
systemctl restart autohotspot.service
